# packaging/app-dropbox.cron
05 2 * * * root /usr/sbin/dropboxconf >/dev/null 2>&1
